# supervisord/filebeat.conf
[program:filebeat]
command=/bin/bash -c "/usr/bin/filebeat -e"
redirect_stderr=true