# supervisord/mysql.conf
[program:mysql]
command=/bin/bash -c "/usr/sbin/mysqld"
redirect_stderr=true
